# supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/app/storage/logs/supervisord.log
pidfile=/var/run/supervisor/supervisord.pid
childlogdir=/app/storage/logs
loglevel=info

[unix_http_server]
file=/var/run/supervisor/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor/supervisor.sock

# PHP-FPM for Roundcube
[program:php-fpm]
command=/usr/sbin/php-fpm8.2 -F
autostart=true
autorestart=true
priority=10
stdout_logfile=/app/storage/logs/php-fpm.log
stderr_logfile=/app/storage/logs/php-fpm.err.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB

# Nginx web server
[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
priority=20
stdout_logfile=/app/storage/logs/nginx.log
stderr_logfile=/app/storage/logs/nginx.err.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB

# Postfix mail server
[program:postfix]
command=/usr/local/bin/billionmail/postfix-wrapper.sh
autostart=true
autorestart=true
priority=30
stdout_logfile=/app/storage/logs/postfix.log
stderr_logfile=/app/storage/logs/postfix.err.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
startsecs=0

# Dovecot IMAP/POP3 server
[program:dovecot]
command=/usr/sbin/dovecot -F
autostart=true
autorestart=true
priority=40
stdout_logfile=/app/storage/logs/dovecot.log
stderr_logfile=/app/storage/logs/dovecot.err.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB

# Rspamd spam filter
[program:rspamd]
command=/usr/bin/rspamd -f
autostart=true
autorestart=true
priority=50
stdout_logfile=/app/storage/logs/rspamd.log
stderr_logfile=/app/storage/logs/rspamd.err.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
user=_rspamd
environment=HOME="/var/lib/rspamd"

# OpenDKIM
[program:opendkim]
command=/usr/sbin/opendkim -f
autostart=true
autorestart=true
priority=60
stdout_logfile=/app/storage/logs/opendkim.log
stderr_logfile=/app/storage/logs/opendkim.err.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB

# ClamAV daemon
[program:clamd]
command=/usr/sbin/clamd --foreground=true
autostart=false
autorestart=true
priority=70
stdout_logfile=/app/storage/logs/clamd.log
stderr_logfile=/app/storage/logs/clamd.err.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
startsecs=60

# ClamAV freshclam (virus database updater)
[program:freshclam]
command=/usr/bin/freshclam -d --foreground=true
autostart=false
autorestart=true
priority=71
stdout_logfile=/app/storage/logs/freshclam.log
stderr_logfile=/app/storage/logs/freshclam.err.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB

# Fail2ban
[program:fail2ban]
command=/usr/bin/fail2ban-server -xf start
autostart=false
autorestart=true
priority=80
stdout_logfile=/app/storage/logs/fail2ban.log
stderr_logfile=/app/storage/logs/fail2ban.err.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB

# BillionMail API server
[program:billionmail-api]
command=/opt/billionmail/venv/bin/gunicorn --bind [::]:8080 --workers 4 --timeout 120 --access-logfile /app/storage/logs/billionmail-api-access.log --error-logfile /app/storage/logs/billionmail-api-error.log app:app
directory=/opt/billionmail
autostart=true
autorestart=true
priority=90
stdout_logfile=/app/storage/logs/billionmail-api.log
stderr_logfile=/app/storage/logs/billionmail-api.err.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
environment=PATH="/opt/billionmail/venv/bin:%(ENV_PATH)s",
            PYTHONPATH="/opt/billionmail:%(ENV_PYTHONPATH)s"

# BillionMail Celery worker for background tasks
[program:billionmail-celery]
command=/opt/billionmail/venv/bin/celery -A app.celery worker --loglevel=info
directory=/opt/billionmail
autostart=true
autorestart=true
priority=91
stdout_logfile=/app/storage/logs/billionmail-celery.log
stderr_logfile=/app/storage/logs/billionmail-celery.err.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
environment=PATH="/opt/billionmail/venv/bin:%(ENV_PATH)s",
            PYTHONPATH="/opt/billionmail:%(ENV_PYTHONPATH)s"

# BillionMail Celery beat for scheduled tasks
[program:billionmail-beat]
command=/opt/billionmail/venv/bin/celery -A app.celery beat --loglevel=info
directory=/opt/billionmail
autostart=true
autorestart=true
priority=92
stdout_logfile=/app/storage/logs/billionmail-beat.log
stderr_logfile=/app/storage/logs/billionmail-beat.err.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
environment=PATH="/opt/billionmail/venv/bin:%(ENV_PATH)s",
            PYTHONPATH="/opt/billionmail:%(ENV_PYTHONPATH)s"

# Groups for managing related services
[group:mail]
programs=postfix,dovecot,rspamd,opendkim
priority=999

[group:web]
programs=nginx,php-fpm
priority=998

[group:billionmail]
programs=billionmail-api,billionmail-celery,billionmail-beat
priority=997

[group:security]
programs=clamd,freshclam,fail2ban
priority=996